# WebGis_Projeto/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn Gis_Project.wsgi:application --bind 0.0.0.0:8000
directory=/code
user=root
autostart=true
autorestart=true
redirect_stderr=true
